# docker/etc/s6-overlay/s6-rc.d/svc-wsebrd/run
#!/usr/bin/with-contenv bash

cd /app && /app/wsebrd --database /data/database.db
